# noloops-daemon-deb/src/DEBIAN/prerm
#!/bin/bash
systemctl stop noloops
systemctl disable noloops
